# package/lrd/lrd-encrypted-storage-toolkit/rootfs/usr/lib/systemd/system/var-log-journal.mount
[Unit]
Description=Persistent Journal Storage Bind
RequiresMountsFor=/perm

[Mount]
What=/perm/log/journal
Where=/var/log/journal
Type=none
Options=bind

[Install]
WantedBy=systemd-journald.service
