# board/wb50n_sysd/rootfs-additions/usr/lib/systemd/system/bluetooth.service.d/05-service-fix.conf
[Service]
ProtectSystem=strict
StateDirectory=/var/lib/bluetooth
ConfigurationDirectory=/etc/bluetooth
ReadWritePaths=
ReadOnlyPaths=
