# package/lrd/lrd-automount/usb-update@.service
[Unit]
BindTo=%i.device
After=%i.device

[Service]
ExecStart=/usr/bin/usb-update.sh /%I
KeyringMode=inherit
