# package/lrd/lrd-encrypted-storage-toolkit/rootfs/usr/lib/systemd/system/nginx-unit.service.d/10-inherit-keyring.conf
[Service]
KeyringMode=inherit

[Unit]
Wants=rodata.mount
After=rodata.mount
